# docker/volume.skel/etc/crontabs/root
MAILTO={{admin_addr}}

0	9	*	*	*	/check.ssl.py > /dev/null 2>&1
